# deluged.conf
[program:deluged]
autorestart=true
command=/usr/bin/deluged --do-not-daemonize
user=deluge
environment=HOME="/home/deluge",USER="deluge"
