# openvpn.conf
[program:openvpn]
command=openvpn --config /etc/openvpn/client.conf
autorestart=true
directory=/etc/openvpn/

