# vnstatd.conf
[program:vnstatd]
command=/usr/sbin/vnstatd --nodaemon
autorestart=true
user=deluge
environment=HOME="/home/deluge",USER="deluge"
